# templates/centos/uhostappserver.conf.erb
#Upstart Service For Nodejs in /etc/init
description 'uhostappserver'
 
env NAME=uhostappserver
env LOG_FILE=/var/log/uhostappserver.log
env USER=uhost
env NODE_BIN=/usr/bin/node
env NODE_ENV="prod"
env PORT=8888
env DIRECTORY=/srv/uhostappserver/current
env APP=app.js
 
start on runlevel [23]
stop on shutdown
# Respawn in case of a crash, with default parameters
respawn
respawn limit 10 5
 
script
  cd $DIRECTORY
  su $USER
  # Make sure logfile exists and can be written by the user we drop privileges to
  touch $LOG_FILE
  chown $USER:$USER $LOG_FILE
  # recommended approach in case of su/sudo usage so that service does not fork
  exec sudo -u $USER PORT=$PORT NODE_ENV=$NODE_ENV $NODE_BIN $DIRECTORY/$APP >> $LOG_FILE 2>&1
end script
 
post-start script
  echo "app $NAME $NODE_ENV post-start event" >> $LOG_FILE
end script
